# server/planary-api.service
[Unit]
Description=Planary API server (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=planary
WorkingDirectory=/opt/planary/app
EnvironmentFile=/opt/planary/app/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Hardening (same as worker)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
